# third_party/m_log/cfg/logd.conf
[program:logd]
command=/usr/local/bin/logd
directory=/tmp
autostart=true
autorestart=true
startretries=3
user=root
stdout_logfile=/var/log/logd.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/logd.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5